# etc/systemd/ollama.service.d/override.conf
[Service]
Environment="OLLAMA_CONTEXT_LENGTH=64000"
Environment="OLLAMA_MODELS=/data/ollama"
